# templates/cloudsshkey.rcd.OpenBSD.j2
#!/bin/sh
#

daemon="/usr/local/sbin/cloudsshkey"

. /etc/rc.d/rc.subr

#pexp="pflogd: \[priv\]"

rc_cmd $1
